# s3.conf
[program:s3]
directory=/opt/src/gcp-imagery-reader/
command=s3fs landsat-pds /imagery -o passwd_file=/opt/src/gcp-imagery-reader/passwd
autostart=true
autorestart=true
redirect_stderr=true
